# deploy/linux/z-cmdb.service
[Unit]
Description=Z-CMDB Asset Management System
After=network.target

[Service]
Type=simple
User=cmdb
Group=cmdb
WorkingDirectory=/opt/z-cmdb/backend
EnvironmentFile=/opt/z-cmdb/backend/.env
ExecStart=/opt/z-cmdb/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
